# pipeline/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
user=root

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Worker ingestion MQTT (continu) ──────────────────────────────────────────
; Reçoit tous les messages MQTT, valide, déduplique et écrit dans InfluxDB raw.
[program:ingestion]
command=python -m workers.ingestion
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/ingestion.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/ingestion.err.log
stderr_logfile_maxbytes=2MB
priority=10

; ── Worker calibration (toutes les 30 s) ─────────────────────────────────────
; Lit bucket_raw, applique Kalman + RF/fallback, écrit bucket_cleansed.
[program:calibration]
command=python -m workers.calibration
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/calibration.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/calibration.err.log
stderr_logfile_maxbytes=2MB
priority=20

; ── Worker détection d'anomalies (toutes les 60 s) ───────────────────────────
; 3 niveaux : seuils IQA → Isolation Forest → règles structurelles.
[program:anomaly]
command=python -m workers.anomaly_detector
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/anomaly.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/anomaly.err.log
stderr_logfile_maxbytes=2MB
priority=20

[group:workers]
programs=metrics,ingestion,calibration,anomaly

; ── Serveur de métriques Prometheus (port 9090) ──────────────────────────────
; Un seul serveur partagé par tous les workers (évite le conflit de port).
[program:metrics]
command=python -c "from metrics import start_metrics_server; import time; s = start_metrics_server(9090); time.sleep(999999)"
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/app/logs/metrics.log
stdout_logfile_maxbytes=2MB
priority=1
